# scripts/launchd/com.cascade-identity.libcoder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cascade-identity.libcoder</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>PERP_ORCHESTRATOR_DIR/dist/cascade-identity.js</string>
        <string>--target</string>
        <string>libcoder</string>
        <string>--workspace</string>
        <string>WORKSPACE_PATH</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SHARED_STATE_DIR</key>
        <string>SHARED_STATE_DIR</string>
        <key>TARGETS_DIR</key>
        <string>TARGETS_DIR</string>
        <key>IDENTITY_HEARTBEAT_MS</key>
        <string>30000</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>LOG_DIR/cascade-identity-libcoder.log</string>

    <key>StandardErrorPath</key>
    <string>LOG_DIR/cascade-identity-libcoder.log</string>
</dict>
</plist>
