# scripts/launchd/com.pine-guard.perp-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pine-guard.perp-orchestrator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/eugene/CascadeProjects/windsurf-project-4/perp-orchestrator/dist/index.js</string>
    <string>--transport</string>
    <string>http</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/eugene/CascadeProjects/windsurf-project-4/perp-orchestrator</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/perp-orchestrator.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/perp-orchestrator.err</string>
</dict>
</plist>
